# debian/postinst
#!/bin/bash
set -e

#if [ "$1" = "configure" ]; then
#    systemctl daemon-reload
#    systemctl enable bind_null_device.service
#    systemctl start bind_null_device.service
#fi

exit 0
