# debian/prerm
#!/bin/bash
set -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    systemctl stop bind_null_device.service
    systemctl disable bind_null_device.service
fi

exit 0
